# eduface.slice
[Unit]
Description=Limited resources Slice for EduFace
DefaultDependencies=no
Before=slices.target

[Slice]
MemoryLimit=12G